# images/browser/supervisor/websockify.conf
; Bridges the live-view WebSocket on :8080 to x11vnc's RFB on localhost:5900.
; This is the port the harness proxies as live_view_url (SERVICE_PORT_LIVE_VIEW
; 443 -> TARGET_PORT_LIVE_VIEW 8080). websockify connects to the RFB target
; lazily per client, so it tolerates x11vnc not being up yet.
[program:websockify]
command=/usr/bin/websockify --heartbeat=30 0.0.0.0:8080 127.0.0.1:5900
autostart=true
autorestart=true
startsecs=2
startretries=100000
stopsignal=TERM
stdout_logfile=/var/log/supervisord/websockify
redirect_stderr=true
