# images/browser/supervisor/x11vnc.conf
; RFB live view: x11vnc attaches to the existing headful Xorg :1 (the same
; display chromium and mutter use) and serves RFB on localhost:5900. websockify
; bridges it to a WebSocket on :8080. No VNC password — the display is bound to
; localhost and access is enforced upstream (control lock + token auth in the
; harness/ops proxy). No -auth: Xorg :1 is started with -nolisten tcp and no
; Xauthority, so local unix-socket clients (neko, mutter, x11vnc) connect freely.
[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -rfbport 5900 -localhost -forever -shared -nopw -noxdamage -norepeat -nolookup -wait 20 -defer 20 -quiet
environment=DISPLAY=":1"
autostart=true
autorestart=true
startsecs=2
startretries=100000
stopsignal=TERM
stdout_logfile=/var/log/supervisord/x11vnc
redirect_stderr=true
